# trunk/crontabs
# minute  hour   day  week  month   command
       0    10     *     *      *   rsync -az /Users/francis/Torrents /Volumes/BigData/.rsync_Torrents
       *     *     *     *      *   activate() { echo "$(date) (30s)" >> /Volumes/BigData/Activate.txt; } ; activate && sleep 30 && activate
